# etc/xdg/autostart/QSyncthingTray.desktop
[Desktop Entry]
Exec=QSyncthingTray
Name=QSyncthingTray
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true