# etc/xdg/autostart/SyncthingTray.desktop
[Desktop Entry]
Exec=syncthingtray
Name=SyncthingTray
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
