# etc/xdg/autostart/user-services.desktop
[Desktop Entry]
Exec=/etc/px/services/start.sh
Name=user-services
OnlyShowIn=LXQt;
Type=Application
Version=1.0