# etc/xdg/autostart/user-services.desktop
[Desktop Entry]
Exec=/bin/px-user-services.sh
Name=user-services
OnlyShowIn=LXQt;
Type=Application
Version=1.0